# apps/server/deploy/huma-x11vnc.service
[Unit]
Description=HUMA x11vnc (Xvfb :99, port 5900)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@HUMA_USER@
Group=@HUMA_USER@
Environment=DISPLAY=:99
Environment=HUMA_VNC_PORT=5900
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=@HUMA_HOME@/huma
ExecStart=/bin/bash @HUMA_HOME@/huma/apps/server/deploy/scripts/run-x11vnc.sh
Restart=always
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
